# outreach/com.dentdigital.crmsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dentdigital.crmsync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/umar/.dentdigital/outreach/crm-sync.sh</string>
  </array>

  <!-- Every 2h during clinic hours (Oslo). Vercel's daily cron is the backstop. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>10</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>10</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>10</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>10</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>10</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>10</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>10</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/umar/.dentdigital/outreach/crmsync.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/umar/.dentdigital/outreach/crmsync.launchd.err.log</string>
</dict>
</plist>
