# outreach/com.dentdigital.outreach.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dentdigital.outreach</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/umar/Documents/DEV/DentDigital/DentDigital/outreach/daily-send.sh</string>
    <string>15</string>
  </array>

  <!-- Weekdays (Mon=1 … Fri=5) at 09:20. B2B-friendly, avoids weekends. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
  </array>

  <!-- If the Mac was asleep at 09:20, run at next wake instead of skipping. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/umar/Documents/DEV/DentDigital/DentDigital/outreach/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/umar/Documents/DEV/DentDigital/DentDigital/outreach/launchd.err.log</string>
</dict>
</plist>
